# packaging/linux/lindyndns.timer
[Unit]
Description=Run Linode Dynamic DNS client every ten minutes

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
